# systemd/monitor.service
[Unit]
Description=Monitor process and restart if not running

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/scripts/monitor.sh
